# deploy/systemd/q-s2-dashboard.service
[Unit]
Description=Q Portal (S2 / S4 / S6 unified)
After=network.target

[Service]
Type=simple
User=ubuntu
WorkingDirectory=/tmp/qs2_review
Environment=PORTAL_HOST=0.0.0.0
Environment=PORTAL_PORT=3010
Environment=S2_RUNTIME_PATH=/tmp/qs2_review
Environment=S2_DB_PATH=/tmp/qs2_review/data/s2.sqlite
Environment=S2_ENV_PATH=/home/ubuntu/.openclaw/.env
# Set PORTAL_PASSWORD to enable auth (leave empty to run open)
Environment=PORTAL_PASSWORD=
# MDX subscription renewal date (YYYY-MM-DD) — colour-coded countdown on /s2
Environment=MDX_RENEWAL_DATE=2026-05-25
# Portfolio baseline at last capital reset — used for P&L and per-bot drawdown criteria
Environment=PORTFOLIO_BASELINE_USDT=2799.94
ExecStart=/usr/bin/env node /tmp/qs2_review/scripts/run-portal.js
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
