# deploy/systemd/q-s2-analysis.service
[Unit]
Description=Q S2 Signal Analysis Processor
After=network.target

[Service]
Type=oneshot
User=ubuntu
Environment=S2_DB_PATH=/home/ubuntu/.openclaw/workspace/Q_S2/data/s2.sqlite
ExecStart=/usr/bin/python3 /tmp/qs2_review/scripts/process_analysis_queue.py
StandardOutput=journal
StandardError=journal
